# crontab
*/30 * * * * python /refresh_playlists.py >> /logfile